# panic-monitor.service
[Unit]
Description=Panic Monitor - P2P Health Monitoring Daemon
After=network-online.target
Wants=network-online.target

# Automatic restart backoff: systemd will stop restarting after 5 failures
# within 30 seconds, preventing a crash loop from consuming resources.
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple

# ── Password ──────────────────────────────────────────────────────────
# The daemon needs the identity password to decrypt secret.key at startup.
# Option A: Set it here (leaks into /proc/<pid>/environ on shared hosts).
# Option B: Use EnvironmentFile= pointing to a root-only-readable file.
#
# Option A (inline):
# Environment=PANIC_MONITOR_PASSWORD=your-password-here
#
# Option B (recommended — create /etc/panic-monitor/env with:
#   PANIC_MONITOR_PASSWORD=your-password-here
#   and chmod 600 /etc/panic-monitor/env):
EnvironmentFile=/etc/panic-monitor/env

ExecStart=/usr/bin/python3 /opt/panic-monitor/main.py --daemon \
    --identity      /etc/panic-monitor/secret.key \
    --identity-meta /etc/panic-monitor/secret.meta \
    --peers         /etc/panic-monitor/peers.json \
    --log-path      /etc/panic-monitor/log.jsonl \
    --history-db    /var/lib/panic-monitor/history.db

WorkingDirectory=/opt/panic-monitor
Restart=on-failure
RestartSec=5

# ── Watchdog ──────────────────────────────────────────────────────────
# systemd will SIGABRT the process if it doesn't notify within 90s.
# (Requires sd_notify integration or Type=notify; with Type=simple this
#  acts as a liveness timeout — if the process hangs hard, systemd kills it.)
WatchdogSec=90

# ── Memory limits ─────────────────────────────────────────────────────
# MemoryHigh = soft limit (triggers kernel memory reclaim pressure).
# MemoryMax  = hard ceiling (OOMKilled if exceeded — prevents runaway leaks
#              from taking down the host).
MemoryHigh=180M
MemoryMax=256M

# ── CPU limits ────────────────────────────────────────────────────────
# A monitoring daemon should never saturate a core. 50% of one core is
# generous (argon2 at startup is the peak; heartbeat cycles are I/O bound).
CPUQuota=50%
# Nice the process so it yields to production workloads.
Nice=10

# ── File descriptor limit ─────────────────────────────────────────────
# Each iroh QUIC connection uses an fd. 512 is generous for a daemon
# monitoring dozens of peers + SQLite + HTTP server.
LimitNOFILE=512

# ── Process limit ─────────────────────────────────────────────────────
# Prevent fork bombs / runaway thread spawning. The daemon needs:
#   main thread + APScheduler thread + status page thread + thread pool.
# 64 is very generous.
LimitNPROC=64

# ── Disk write rate limit (optional) ──────────────────────────────────
# Prevent a bug that hammers SQLite from saturating disk I/O.
# 10 MB/s is generous for a monitoring daemon that writes tiny probe rows.
IOWriteBandwidthMax=/var/lib/panic-monitor 10M

# ── Hardening ─────────────────────────────────────────────────────────
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
ReadWritePaths=/etc/panic-monitor /var/lib/panic-monitor

# Restrict network to only AF_INET/AF_INET6/AF_UNIX (no raw sockets, netlink, etc.)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Restrict syscalls to a safe set (no personality, kexec, etc.)
SystemCallFilter=@system-service
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
