# panic-monitor.service
[Unit]
Description=Panic Monitor - P2P Health Monitoring Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/panic-monitor/main.py --daemon --peers /etc/panic-monitor/peers.json --identity /etc/panic-monitor/secret.key
WorkingDirectory=/opt/panic-monitor
Restart=on-failure
RestartSec=5

# Memory limits
MemoryMax=200M
MemoryHigh=180M

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/panic-monitor

[Install]
WantedBy=multi-user.target
